# supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock

[supervisord]
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=50MB
logfile_backups=5
loglevel=info
nodaemon=true
pidfile=/var/run/supervisord.pid

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

# ── Node.js Express + Socket.IO API ───────────────────────────────────────────
[program:nodejsbackend]
command=node /app/server.js
directory=/app
environment=NODE_ENV="production",PYTHON_SERVICE_URL="http://127.0.0.1:5002"
autostart=true
autorestart=true
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

# ── FastAPI Python service (Uvicorn) ──────────────────────────────────────────
[program:pythonservice]
command=uvicorn app:app --host 0.0.0.0 --port 5002 --workers 1
directory=/app/python_service
autostart=true
autorestart=true
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=5
